# templates/certbot.cron
12 3 * * 0 root certbot renew --standalone && systemctl reload prosody
